# redis/files/supervisor.conf
{%- from "linux/map.jinja" import network with context %}
{%- set service = salt['pillar.get']('supervisor:server:service:'+service_name) %}

[program:{{ service.type }}_{{ service.name }}]
command=redis-sentinel /etc/redis/sentinel.conf
stdout_logfile=/var/log/redis/redis-sentinel.out.log
stderr_logfile=/var/log/redis/redis-sentinel.err.log
autostart=true
autorestart=true

